# storage/framework/supervisor/prompt-flow-worker.conf
[program:prompt-flow-worker]
process_name=%(program_name)s_%(process_num)02d
command=/usr/bin/php8.4 /home/bruno/Documents/projects/b7s/ai-project/artisan queue:work --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
user=bruno
numprocs=2
redirect_stderr=true
redirect_stdout=true
stdout_logfile=/home/bruno/Documents/projects/b7s/ai-project/storage/logs/worker.log
stderr_logfile=/home/bruno/Documents/projects/b7s/ai-project/storage/logs/worker-error.log
stopwaitsecs=3600